# facepalm.service
[Unit]
Description=Facepalm

[Service]
ExecStart=/usr/bin/facepalm
Restart=always

[Install]
WantedBy=multi-user.target